# deploy/pysis.supervisor.conf
[program:pysis]
command=/virtualenvs/pysis/bin/python /projects/pysis/deploy/pysis_fcgi.py
user=www-data
